# updates_check.timer
[Unit]
Description=Run updates_check script periodically

[Timer]
OnBootSec=2min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
