# updates_check.service
[Unit]
Description=Package Updates Count Exporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/updates_check.sh
